# templates/radarr.service.erb
############################################
# WARNING:
# Managed by puppet - DO NOT modify locally.
############################################

[Unit]
Description=Radarr Daemon
After=syslog.target network.target

[Service]
User=<%= scope['::radarr::user'] %>
Group=<%= scope['::radarr::group'] %>
Type=simple
Restart=on-failure
#TODO parameterise path
ExecStart=/usr/bin/mono <%= scope['::radarr::install_path'] %>/Radarr/Radarr.exe --nobrowser
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target

