# templates/radarr.service.erb
############################################
# WARNING:
# Managed by puppet - DO NOT modify locally.
############################################

[Unit]
Description=Radarr Daemon
After=network.target

[Service]
User=<%= scope['::radarr::user'] %>
Group=<%= scope['::radarr::group'] %>
Type=simple
Restart=always
RestartSec=5
#TODO parameterise path
ExecStart=/usr/bin/mono --debug <%= scope['::radarr::install_path'] %>/Radarr/RadarrConsole.exe --NoRestart
TimeoutStopSec=20
#ExecStart=/usr/bin/radarr-daemon --foreground --config-dir <%= scope['::radarr::params::config_dir'] %><% if scope['::radarr::log_file'] -%> --logfile <%= scope['::radarr::log_file'] %><% end %>
#ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target

